# deploy/vtt-remote.service
# VTT Remote Relay - systemd service
#
# Installation:
#   1. Build: task build:server
#   2. Copy binary to /opt/vtt-remote/vtt-relay
#   3. Create user: useradd -r -s /bin/false vttremote
#   4. Set ownership: chown -R vttremote:vttremote /opt/vtt-remote
#   5. Copy this file to /etc/systemd/system/vtt-remote.service
#   6. sudo systemctl daemon-reload
#   7. sudo systemctl enable --now vtt-remote

[Unit]
Description=VTT Remote Relay Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=vttremote
Group=vttremote

# Working directory
WorkingDirectory=/opt/vtt-remote

# Binary location
ExecStart=/opt/vtt-remote/vtt-relay

# Environment (server defaults to 8080)
Environment=PORT=8080

# Restart policy
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true

# Resource limits
LimitNOFILE=65535
MemoryMax=256M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vtt-remote

[Install]
WantedBy=multi-user.target
